# deploy/systemd/jcv-collect-daily-events.service
[Unit]
Description=JCV cron · collect-daily-events
After=network.target

[Service]
Type=oneshot
User=chenkuichen
WorkingDirectory=/home/chenkuichen/app/portfolio-dashboard
EnvironmentFile=/home/chenkuichen/app/portfolio-dashboard/.env.local
ExecStart=/usr/bin/curl -fsSL --max-time 300 \
  -H "Authorization: Bearer ${CRON_SECRET}" \
  http://127.0.0.1:3000/api/cron/collect-daily-events
StandardOutput=journal
StandardError=journal
