# deploy/systemd/jcv-collect-target-prices.timer
[Unit]
Description=JCV cron timer · collect-target-prices

[Timer]
OnCalendar=Mon..Fri *-*-* 09:30:00
Persistent=true
Unit=jcv-collect-target-prices.service

[Install]
WantedBy=timers.target
