# deploy/systemd/jcv-compute-valuations.timer
[Unit]
Description=JCV cron timer · compute-valuations

[Timer]
OnCalendar=*-*-* 22:00:00 UTC
Persistent=false
AccuracySec=1m
Unit=jcv-compute-valuations.service

[Install]
WantedBy=timers.target
