# deploy/systemd/jcv-collect-news.timer
[Unit]
Description=JCV cron timer · collect-news

[Timer]
OnCalendar=Mon..Fri *-*-* 10:00:00 UTC
Persistent=false
AccuracySec=1m
Unit=jcv-collect-news.service

[Install]
WantedBy=timers.target
