# deploy/systemd/jcv-collect-daily-events.timer
[Unit]
Description=JCV cron timer · collect-daily-events

[Timer]
OnCalendar=Mon..Fri *-*-* 08:00:00
Persistent=true
Unit=jcv-collect-daily-events.service

[Install]
WantedBy=timers.target
